# crontab
*/10 * * * * /bin/sh /ping_script.sh
5 */1 * * * /usr/bin/speedtest-cli > /config/speedtest
# An empty line is required at the end of this file for a valid cron file.
